# pithagoras.service
[Unit]
Description=Pithagoras

[Service]
ExecStart=/usr/local/bin/pithagoras
Restart=always

[Install]
WantedBy=multi-user.target